# systemd/user/sway.service
[Unit]
Description=Sway WM

[Service]
ExecStart=/usr/bin/sway

[Install]
WantedBy=default.target
